# systemd/openqa_bugfetcher@.timer
[Unit]
Description=openQA bug fetching task

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
